# scrapy-samples/estate/estate.conf
# supervisor conf


[program:163_bj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:163_cq]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:163_gz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:163_sh]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:163_sz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:163_tj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:163_xm]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:cei]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:chinanews]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:fang]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:funxun_bj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:funxun_gz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:funxun_sh]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:funxun_sz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:funxun_tj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:haofung]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:ifeng_bj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:ifeng_cq]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:ifeng_gz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:ifeng_sh]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:ifeng_sz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:ifeng_tj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:payfang_bj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:payfang_gz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:payfang_sh]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:payfang_sz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:payfang_tj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:people]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:qq_bj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:qq_cq]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:qq_gz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:qq_sh]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:qq_sz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:qq_tj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:qq_xm]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sina_bj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sina_cq]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sina_gz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sina_sh]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sina_sz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sina_tj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sina_xm]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sohu_cq]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sohu_gz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sohu_sh]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sohu_sz]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sohu_tj]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:sohu_xm]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log


[program:thepaper]
command=/home/datartisan/works/python/envs/scrapy/bin/scrapy crawl %(program_name)s
process_name=%(program_name)s
directory=/home/datartisan/works/python/envs/scrapy/scrapy-estate
user=datartisan
autostart=false
stdout_logfile=/home/datartisan/works/log/estate-%(program_name)s-log
stderr_logfile=/home/datartisan/works/log/estate-error-%(program_name)s-log
